# deploy/dod-ocr-backend.service
[Unit]
Description=DoD OCR FastAPI backend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/dod-ocr/backend
EnvironmentFile=/opt/dod-ocr/backend/.env
ExecStart=/home/ubuntu/.local/bin/uv run uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
